# wangemu-terminal-server.service
[Unit]
Description=Wang 2200 Terminal Server Emulator
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/wangemu
ExecStart=/home/pi/wangemu/wangemu-terminal-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target